# trafilatura-extractor.service
[Unit]
Description=Trafilatura Extractor FastAPI Application
After=network.target

[Service]
Type=simple
User=_SERVICE_USER_
WorkingDirectory=/home/_SERVICE_USER_/trafilatura-extractor
ExecStart=/home/_SERVICE_USER_/trafilatura-extractor/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/trafilatura-extractor/app.log
StandardError=append:/var/log/trafilatura-extractor/error.log
EnvironmentFile=-/home/_SERVICE_USER_/trafilatura-extractor/.env

[Install]
WantedBy=multi-user.target
